# root/etc/s6-overlay/s6-rc.d/nginx-redirect-conf/up
with-contenv
foreground { /etc/s6-overlay/s6-rc.d/nginx-redirect-conf/redirect-conf.sh }
